# helios-raze/helios-seed-photonvision-camera.service
[Unit]
Description=Seed PhotonVision camera configuration
After=local-fs.target
Before=photonvision.service
ConditionPathExists=/opt/photonvision/photonvision_config

[Service]
Type=oneshot
ExecStart=/usr/local/bin/helios-seed-photonvision-camera.sh

[Install]
WantedBy=multi-user.target
